# src/init/openrc
#!/sbin/openrc-run
# Copyright (c) 2024 TheRealOne78
# Distributed under the terms of the GNU General Public License v3+
# https://www.gnu.org/licenses/gpl-3.0.en.html

name="ntfy-resources-usage-plugin"
description="Send NTFY alerts to your smart-phone when your unix computer is running out of resources."

supervisor="supervise-daemon"
command="/bin/bash /usr/bin/ntfy-resources-usage"
command_background="true"
